# systemd/rsync-backup.timer
[Unit]
Description=Backup with rsync weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
